# templates/gentoo/initramfs/etc/local.d/sshd.start
#!/bin/sh
rc-service net.* start || true
ssh-keygen -A -f /etc
rc-service sshd restart || true
